# deploy/systemd/gadgetd-control.service
# TeslaUSB B-1 — gadgetd-control: the eject-handoff control daemon. It serves
# the mutator IPC over a Unix socket (the ONLY Pi-side write authority for the
# car-facing image, SPEC.md §4) and performs handoff crash-recovery on start.
#
# DECOUPLED FROM BRING-UP BY DESIGN: gadget bring-up is owned by gadgetd.service
# (oneshot). This daemon NEVER tears the gadget down — it has no gadget-touching
# ExecStop — so a crash, restart, or socket failure here can never disturb the
# car-facing LUN. It only ejects/re-presents the medium for the bounded window
# of a serialized handoff.
#
# Install: copy to /etc/systemd/system/gadgetd-control.service, then
#   systemctl daemon-reload && systemctl enable --now gadgetd-control.service
[Unit]
Description=TeslaUSB gadgetd control daemon (eject-handoff mutator IPC)
Documentation=https://github.com/mphacker/TeslaUSB
# The gadget must already be up; we serve control for it, never bring it up.
Requires=gadgetd.service
After=gadgetd.service

[Service]
Type=exec
# Guardian of the write path: never an OOM-kill victim.
OOMScoreAdjust=-1000
# NOTE: --allow-hot-handoff is intentionally OMITTED. Production handoffs are
# refused while the host is enumerated (car connected) until prototype-unknown
# #2 — the car's eject tolerance during an active save — is HW-measured
# (SPEC.md §9). Enable the flag only for bench validation, never on the car.
ExecStart=/usr/local/bin/gadgetd serve --image /data/teslausb/teslacam.img --media-image /data/teslausb/media.img
# Create /run/teslausb (socket + per-handoff working dirs), root-owned, mode 0750.
# PRESERVE: /run/teslausb is SHARED with scannerd.service (both bind sockets
# here). systemd does NOT reference-count RuntimeDirectory across units, so
# without this a `systemctl restart scannerd` would tear down /run/teslausb and
# unlink THIS daemon's live gadgetd.sock — leaving gadgetd LISTENing on an
# orphaned inode while webd's connect() fails with ENOENT (the write path dies
# silently). Preserve=yes keeps the dir across any peer service stop/restart.
RuntimeDirectory=teslausb
RuntimeDirectoryMode=0750
RuntimeDirectoryPreserve=yes
# A serve/socket failure must not be fatal to the system; just restart serving.
Restart=always
RestartSec=2
# There is NO ExecStop=down here, on purpose: stopping the control daemon must
# leave the gadget bound and the medium present.

[Install]
WantedBy=multi-user.target
