# deploy/systemd/gadgetd.service
# TeslaUSB B-1 — gadgetd: the kernel mass-storage gadget that presents the
# car-facing USB drive. configfs is tmpfs and is lost on every reboot, so this
# unit rebuilds the gadget at boot. It is the guardian of the #1 invariant
# (SPEC.md §2): the car must always be able to write TeslaCam.
#
# Install: copy to /etc/systemd/system/gadgetd.service, then
#   systemctl daemon-reload && systemctl enable --now gadgetd.service
[Unit]
Description=TeslaUSB gadgetd (kernel USB mass-storage gadget, car-facing write path)
Documentation=https://github.com/mphacker/TeslaUSB
# configfs must be mounted before we can build the gadget tree.
Requires=sys-kernel-config.mount
After=sys-kernel-config.mount local-fs.target
# A UDC must exist (dwc2 peripheral). If absent, the service fails loudly.
ConditionPathExists=/sys/class/udc

[Service]
Type=oneshot
RemainAfterExit=yes
# The single most important process on the device: never an OOM-kill victim.
OOMScoreAdjust=-1000
ExecStartPre=/usr/local/bin/gadgetd provision --image /data/teslausb/disk.img --size-mib 4096 --media-mib 1024
ExecStart=/usr/local/bin/gadgetd up --image /data/teslausb/disk.img
ExecStop=/usr/local/bin/gadgetd down --image /data/teslausb/disk.img

[Install]
WantedBy=multi-user.target
