# deploy/systemd/indexd.service
# TeslaUSB B-1 — indexd: derives trips/events/clips from the car footage and
# persists them into the SQLite catalog that webd/uploadd/retentiond read.
#
# indexd is the sole DB writer and a pure CONSUMER of scannerd's facts (ADR
# 0001): it NEVER opens the backing image and never parses car-written bytes.
# It connects to scannerd's control socket (/run/teslausb/scannerd.sock),
# drives the 30 s scan cadence, receives a validated batch of facts per pass,
# and applies it to the catalog. This keeps the untrusted-byte parsing in the
# disposable scannerd process and out of the DB-owning one (the
# privilege/fault-isolation seam). It never writes the car-facing LUN (the #1
# invariant, SPEC.md §2).
#
# Hardening note: cgroup MemoryMax, a dedicated non-root User=, and
# ProtectSystem/sandboxing are intentionally NOT set here. They are device-
# calibrated in migration M5 / Task 7.4c (MemoryMax per the microSD-contention
# spike 2.7). The canonical OOM order (SPEC.md §7) IS encoded now via
# OOMScoreAdjust: indexd is the least-disposable of the non-gadgetd services
# (it owns the catalog), so it carries the lowest score above gadgetd's -1000.
#
# Install: copy to /etc/systemd/system/indexd.service, then
#   systemctl daemon-reload && systemctl enable indexd.service
[Unit]
Description=TeslaUSB indexd (footage catalog: trips/events/clips -> SQLite; IPC consumer)
Documentation=https://github.com/mphacker/TeslaUSB
# Needs /var/lib (catalog DB) mounted and the scannerd producer running:
# indexd no longer reads /data itself, so it depends on scannerd (which owns
# the read-only image fd) rather than on the backing image directly.
RequiresMountsFor=/var/lib/teslausb
After=local-fs.target scannerd.service
Requires=scannerd.service

[Service]
Type=simple
# INDEXD_DB = catalog SQLite path (webd's WEBD_DB MUST match this).
# INDEXD_SCANNERD_SOCKET = the scannerd control socket indexd consumes facts
# from. INDEXD_IMAGE is GONE: indexd never opens the backing image (ADR 0001).
Environment=INDEXD_DB=/var/lib/teslausb/index.sqlite3
Environment=INDEXD_SCANNERD_SOCKET=/run/teslausb/scannerd.sock
ExecStart=/usr/local/bin/indexd
Restart=on-failure
RestartSec=5s
# Least-disposable non-gadgetd service: survives memory pressure longest.
OOMScoreAdjust=-100
SyslogIdentifier=indexd

[Install]
WantedBy=multi-user.target
