# deploy/systemd/retentiond.service
# TeslaUSB B-1 — retentiond runs independently of indexd/scannerd daemons:
# it reads teslacam.img directly via linked scannerd library APIs and self-retries
# on startup races where the archive/image appears after boot. I/O is explicitly
# idle/low weight so clip archiving yields to foreground workloads and car-write
# pressure.
[Unit]
Description=TeslaUSB retentiond (archiving + retention + space governor)
Documentation=https://github.com/mphacker/TeslaUSB
RequiresMountsFor=/data/teslausb
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/usr/local/bin/retentiond serve --archive-recent-only --no-delete --archive-root /data/teslausb/archive --volume-image /data/teslausb/teslacam.img
Restart=always
RestartSec=5s
# Catch hangs that do not exit; Restart=always only recovers after failure/exit.
# This is a systemd *service* watchdog (sd_notify WATCHDOG=1), not the hardware
# /dev/watchdog timer. 240s clears the worst-case single-clip copy+hash+fsync
# stall under microSD contention with margin; a false kill is bounded (no footage
# loss — copies are non-destructive staged-promote, markers durable per ADR-0006).
WatchdogSec=240
NotifyAccess=main
IOSchedulingClass=idle
CPUWeight=10
IOWeight=10
OOMScoreAdjust=100
# panic=abort means any panic is a hard exit; suppress coredumps so a
# transient crash-restart cycle can never spam the microSD with dumps.
LimitCORE=0
SyslogIdentifier=retentiond

[Install]
WantedBy=multi-user.target
