# rust/crates/teslausb-worker/units/teslausb-worker.service
# systemd unit for the `teslausb-worker` background daemon.
#
# Unlike `teslafat@.service`, the worker is NOT instanced — there
# is one worker per host. It indexes new dashcam/sentry clips
# (SEI → GPS waypoints in SQLite) and reaps no-GPS RecentClips
# clips per the configured retention + free-space policy.
#
# The worker runs AFTER the FAT/exFAT NBD daemon (`teslafat@0`)
# has bound its socket and the host's `nbd-client` has attached
# the loopback device, because the worker walks the backing
# tree the daemon owns. systemd's `After=` does NOT pull the
# dependency in, only orders it; the actual `Requires=` lives
# on the per-LUN nbd-attach service shipped in Phase 6.x.
#
# Hardening rationale (matches docs/03-CODE-QUALITY-CHARTER.md
# §"Best Architecture Practices" — minimum viable privilege):
#   - The worker opens no network sockets. It reads/writes
#     `/var/lib/teslausb/` (its SQLite store) and reads/deletes
#     under `backing_root` (default `/srv/teslausb`). Anything
#     beyond that is locked down.
#   - SIGTERM is honoured; the supervisor's tokio::select! loop
#     drops the watcher channel, aborts the inotify thread, and
#     exits zero so `systemctl stop` is clean.

[Unit]
Description=TeslaUSB B-1 background worker (indexer + cleanup)
Documentation=https://github.com/mphacker/TeslaUSB
After=teslafat@0.service local-fs.target
Wants=teslafat@0.service

[Service]
ExecStartPre=/usr/local/bin/teslausb-worker --config /etc/teslausb/worker.toml --check-config
ExecStart=/usr/local/bin/teslausb-worker --config /etc/teslausb/worker.toml

User=teslausb
Group=teslausb

CapabilityBoundingSet=
AmbientCapabilities=
NoNewPrivileges=yes

ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/var/lib/teslausb /srv/teslausb

PrivateNetwork=yes
RestrictAddressFamilies=

ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Restart policy. The supervisor's ShutdownReason::is_fatal
# returns false for SIGTERM/SIGINT and BootstrapOnly (exit 0)
# and true for WatcherFailed/IndexerBootstrapFailed/
# CleanupFailed (exit 1). `on-failure` therefore restarts the
# subsystem errors but leaves a clean stop alone.
Restart=on-failure
RestartSec=10s

# `nice` and `ionice` per the worker's "background, never
# disrupt recording" remit. Matches the v1 archive worker's
# behaviour.
Nice=10
IOSchedulingClass=idle

Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
