# templates/chime_scheduler.service
[Unit]
Description=TeslaUSB Chime Scheduler Check
After=present_usb_on_boot.service

[Service]
Type=oneshot
User=__TARGET_USER__
WorkingDirectory=__GADGET_DIR__
ExecStart=/usr/bin/python3 __GADGET_DIR__/scripts/check_chime_schedule.py
StandardOutput=journal
StandardError=journal

# Nice level for low priority
Nice=10

[Install]
WantedBy=multi-user.target
