# templates/cloud_archive_sync.timer
[Unit]
Description=TeslaUSB Cloud Archive Safety Sync
Requires=cloud_archive_sync.service

[Timer]
# Safety fallback — primary trigger is WiFi connect (NetworkManager dispatcher)
# This just ensures sync runs at least once every 24 hours if WiFi stays connected
OnBootSec=10min
OnUnitActiveSec=24h
Persistent=true

[Install]
WantedBy=timers.target
