# templates/cloud_archive_sync.service
[Unit]
Description=TeslaUSB Cloud Archive Sync
After=network-online.target gadget_web.service teslausb-safe-mode.service
Wants=network-online.target
ConditionPathExists=!/run/teslausb-safe-mode

[Service]
Type=oneshot
User=root
# Trigger sync via the web API so the upload thread runs inside the
# long-lived gadget_web process.  Previous approach imported Python
# directly, starting a daemon thread that died when this oneshot exited.
ExecStart=/usr/bin/curl -sf -X POST -H "Content-Type: application/json" -d "{}" http://localhost/cloud/api/sync_now
TimeoutStartSec=30
Nice=15

StandardOutput=journal
StandardError=journal

Restart=on-failure
RestartSec=60
StartLimitBurst=3
StartLimitIntervalSec=600

[Install]
WantedBy=multi-user.target
