# templates/network-optimizations.service
[Unit]
Description=TeslaUSB Network Performance Optimizations
After=network-online.target NetworkManager.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=__GADGET_DIR__/scripts/optimize_network.sh

[Install]
WantedBy=multi-user.target
