# templates/present_usb_on_boot.service
[Unit]
Description=Present USB gadget at boot (with optional cleanup)
After=multi-user.target
Wants=multi-user.target

[Service]
Type=oneshot
ExecStart=__GADGET_DIR__/scripts/boot_present_with_cleanup.sh
RemainAfterExit=yes
User=root

# Restart policy
Restart=on-failure
RestartSec=15
StartLimitBurst=3
StartLimitIntervalSec=600

# Memory limit
MemoryMax=100M

[Install]
WantedBy=multi-user.target