# templates/present_usb_on_boot.service
[Unit]
Description=Present USB gadget at boot
After=teslausb-safe-mode.service
Wants=multi-user.target
# Skip in safe mode (3+ reboots in 10 min)
ConditionPathExists=!/run/teslausb-safe-mode
StartLimitBurst=3
StartLimitIntervalSec=600

[Service]
Type=oneshot
# Present USB FIRST — no cleanup wrapper. Deferred tasks run separately.
ExecStart=/bin/bash __GADGET_DIR__/scripts/present_usb.sh --boot
RemainAfterExit=yes
User=root

# Restart policy
Restart=on-failure
RestartSec=15

# Memory limit
MemoryMax=100M

[Install]
WantedBy=multi-user.target
