# templates/wifi-powersave-off.service
[Unit]
Description=Disable WiFi Power Management and Optimize Performance
After=network.target
Before=network-online.target

[Service]
Type=oneshot
# Disable power management
ExecStart=/sbin/iwconfig wlan0 power off
# Disable roaming (prevent channel switching)
ExecStart=/bin/sh -c 'echo 1 > /sys/module/brcmfmac/parameters/roamoff 2>/dev/null || true'
# Set regulatory domain to US (max power)
ExecStart=/sbin/iw reg set US
# Set TX power to maximum (30 dBm = 1000 mW)
ExecStart=/sbin/iw dev wlan0 set txpower fixed 3000
RemainAfterExit=yes

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
