# templates/wifi-monitor.service
[Unit]
Description=WiFi Connection Monitor and Auto-Reconnect
After=network-online.target teslausb-safe-mode.service
Wants=network-online.target
# WiFi monitor runs even in safe mode so AP fallback works for recovery

[Service]
Type=simple
ExecStart=/bin/bash __GADGET_DIR__/scripts/wifi-monitor.sh
Restart=on-failure
RestartSec=30
StartLimitBurst=5
StartLimitIntervalSec=300
# Allow the script up to 30s to clean up. Stop sequences can be in the middle of
# stop_ap (multi-second teardown) or an nmcli rescan poll (25s) when systemctl
# stop is invoked. With the prior 5s timeout systemd SIGKILL'd mid-cleanup,
# leaving a stale /var/run/wifi-monitor.lock that blocked the next start.
TimeoutStopSec=30
KillMode=mixed

# Memory limit (lightweight script)
MemoryMax=50M

# Run with low priority
Nice=10

[Install]
WantedBy=multi-user.target
